# data/supervisor.conf
[program:mihome]
command = python3 run-mihome.py
directory = /opt/loxberry/data/plugins/mihome/scripts
user = loxberry
environment=HOME="/loxberry", USER="loxberry"
